# cronjobs.cron
2 * * * * curl http://localhost:8082/tasks/select_batch
5 * * * * curl http://localhost:8082/tasks/scrape_views
